# crontab
* * * * * root /letsencrypt/letsencrypt-auto renew >> /var/log/lep/renew.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
